# deploy/grainbidcalc-worker.service
[Unit]
Description=GrainBidCalc - Celery Worker
After=network-online.target redis.service
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/grainbidcalc
ExecStart=/opt/grainbidcalc/venv/bin/celery -A workers.celery_app worker --loglevel=info
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
